# systemd/oturpe-backup.service
[Unit]
Description=Backup personal files

[Service]
Type=simple
WorkingDirectory=/home/username/backup
ExecStart=/home/username/backup/backup
